# deployment/ors-server.service
[Unit]
Description=Open Resource Server in Uzbekistan with ASP.NET Core
After=network.target

[Service]
WorkingDirectory=/path/to/publish/directory/
ExecStart=/usr/bin/dotnet /path/to/publish/directory/ORS.dll
Restart=always
User=www-data
Group=www-data
Environment=ASPNETCORE_URLS=http://localhost:4001 ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target
